# miniredis/deploy/miniredis.service
[Unit]
Description=MiniRedis key-value server
After=network.target

[Service]
Type=simple
User=miniredis
Group=miniredis
WorkingDirectory=/opt/miniredis
Environment=MINIREDIS_BIND=127.0.0.1
Environment=MINIREDIS_PORT=6366
Environment=MINIREDIS_STATS_BIND=127.0.0.1
Environment=MINIREDIS_STATS_PORT=8080
Environment=MINIREDIS_SNAPSHOT_FILE=/var/lib/miniredis/snapshot_6366.dat
Environment=MINIREDIS_SNAPSHOT_INTERVAL=20
EnvironmentFile=-/etc/miniredis/miniredis.env
ExecStart=/opt/miniredis/miniredis
Restart=on-failure
RestartSec=2
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/miniredis

[Install]
WantedBy=multi-user.target
