# deploy/miniredis.service
[Unit]
Description=MiniRedis key-value server
After=network.target

[Service]
Type=simple
User=miniredis
Group=miniredis
WorkingDirectory=/opt/miniredis
EnvironmentFile=-/etc/miniredis/miniredis.env
ExecStart=/opt/miniredis/miniredis --config /etc/miniredis/miniredis.conf
Restart=on-failure
RestartSec=2
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/miniredis /var/log/miniredis

[Install]
WantedBy=multi-user.target
